# ansible/roles/node_app/templates/node-app.service.j2
[Unit]
Description=My Node.js App
After=network.target
[Service]
User=deploy
WorkingDirectory=/home/deploy/apps/myapp
ExecStart=/usr/bin/node /home/deploy/apps/myapp/index.js
Restart=always
Environment=NODE_ENV=production
# Si usas variables secretas, inyectarlas vía environment file o systemd drop-in
[Install]
WantedBy=multi-user.target
Ajusta ExecStart.